# etc/ubuntuhome/.config/upstart/moumon.conf
start on desktop-start
stop on desktop-end

respawn

exec /storage/ca-st-ro/bin/moumon &> /storage/ca-st-ro/moumon.log
